# deploy/autostart/vending-browser.desktop
[Desktop Entry]
Type=Application
Name=Vending Browser
Exec=/home/rpi/bin/start-vending-browser.sh
X-GNOME-Autostart-enabled=true
